# chronos.service
[Unit]
Description=Chronos Main Stack
After=docker.service network.target
Requires=docker.service

[Service]
Type=simple
WorkingDirectory=/opt/chronos
EnvironmentFile=/opt/chronos/.env.deployment

# Ensure docker-compose is down before starting
ExecStartPre=-/usr/bin/docker compose -f docker-compose.yml -f docker-compose.deployment.yml down
ExecStart=/usr/bin/docker compose -f docker-compose.yml -f docker-compose.deployment.yml up --build
ExecStop=/usr/bin/docker compose -f docker-compose.yml -f docker-compose.deployment.yml down

Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
