# edge_server/chronos.service
[Unit]
Description=Chronos Edge Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/chronos
Environment=PATH=/opt/chronos/.venv/bin:$PATH
ExecStart=/opt/chronos/.venv/bin/uvicorn chronos.app:app --host 0.0.0.0 --port 5171
Restart=always
RestartSec=1

[Install]
WantedBy=multi-user.target
